# crontab
# Run ecs_service_discovery.py every 3 minutes
*/3 * * * * python /app/ecs_service_discovery.py >> /var/log/cron.log 2>&1

# An empty line is required at the end of this file for a valid cron file.
